# systemd/megadisplayd.service
[Unit]
Description=MegaDisplay Daemon — USB display for Android tablets
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/megadisplayd
Restart=on-failure
RestartSec=3
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
